# oeasc-service.conf
[program:oeasc]
command = {{APP_PATH}}/flask_run.sh
autostart=true
autorestart=true
stdout_logfile = {{APP_PATH}}/var/log/gn_errors.log
redirect_stderr = true
user = {{USER}}
